# h5serv/main.conf
[program:h5serv]
user=user
directory=/home/user/h5serv/server
command=python app.py
autostart = false
autorestart = false
redirect_stderr = true
