# h5serv/main.conf
[program:h5serv]
user=user
directory=/home/user/h5serv/server
command=python app.py
autostart = true
autorestart = true
redirect_stderr = true
